# rust-client/systemd/epn-client.service
[Unit]
Description=EPN Client Kiosk Application
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/epn-gui
Restart=always
RestartSec=2
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U

# Optionnel: logs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
