# client/systemd/poste-client.service
[Unit]
Description=Client Poste Public - Gestion des sessions
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=poste
Group=poste
WorkingDirectory=/opt/poste-client
ExecStart=/usr/bin/python3 /opt/poste-client/poste_client.py --interactive
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Variables d'environnement
Environment="POSTE_SERVER_URL=http://192.168.1.10:8001"
Environment="POSTE_WS_URL=ws://192.168.1.10:8001"
Environment="LOG_LEVEL=INFO"

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/poste-client.log

[Install]
WantedBy=multi-user.target
