# deployment/ansible/epn-client.service.j2
[Unit]
Description=EPN Client - Gestion Poste Public ({{ inventory_hostname }})
Documentation=https://github.com/mairie-reunion/epn-client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=epn
Group=epn
WorkingDirectory=/opt/epn-client

Environment="RUST_LOG=info"
Environment="EPN_CONFIG=/etc/epn-client/config.yaml"

ExecStart=/usr/local/bin/epn-client

Restart=always
RestartSec=10
StartLimitInterval=0

ExecStartPre=/bin/sleep 5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/epn-client

StandardOutput=journal
StandardError=journal
SyslogIdentifier=epn-client

[Install]
WantedBy=multi-user.target
